# ami-files/etc-systemd-system/npm-preflight.service
[Unit]
Description=Northstar NPM preflight checks (first boot)
Documentation=https://docs.northstar.cloud/npm-ami
Wants=network-online.target
After=network-online.target docker.service
Before=npm-init.service
ConditionPathExists=!/var/lib/npm-init-complete

[Service]
Type=oneshot
ExecStart=/usr/local/bin/npm-preflight
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

