# ami-files/etc-systemd-system/npm-init.service
[Unit]
Description=Nginx Proxy Manager First-Boot Initialization
Documentation=https://docs.northstar.cloud/npm-ami
After=npm.service npm-preflight.service
Requires=npm.service npm-preflight.service
ConditionPathExists=!/var/lib/npm-init-complete

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/npm-init.py
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60
StartLimitBurst=5

[Install]
WantedBy=multi-user.target

